# plans/kmonad/templates/kmonad.service.erb
[Unit]
Description=Launch KMonad

[Service]
Type=simple
StandardOutput=journal
<% if @profile.second? %>
ExecStart=/bin/bash -c "%h/.local/bin/kmonad %h/.config/kmonad/thinkpad-t460.kbd"
<% elsif @profile.third? %>
ExecStart=/bin/bash -c "%h/.local/bin/kmonad %h/.config/kmonad/thinkpad-p15v.kbd"
<% end %>

[Install]
WantedBy=default.target
