# stow/systemd/.config/systemd/user/xset.service
[Unit]
Description=Set caret repeat rate
After=kmonad.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/xset r rate 250 44
